# synapse/rootfs/etc/services.d/synapse/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

export LOG_LEVEL="trace"

CONFIG_DATA_DIR="$(bashio::config 'data_dir')"
DATA_DIR="${CONFIG_DATA_DIR:-"/config"}"
CONFIG_DIR="${DATA_DIR}/config"
SERVER_CONFIG_FILE="${CONFIG_DIR}/homeserver.yaml"

bashio::log.info "Starting Synapse..."

exec /usr/bin/synapse_homeserver --config-path "${SERVER_CONFIG_FILE}"
